# scripts/nfr.service
[Unit]
Description=Network Flight Recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m nfr
WorkingDirectory=/opt/nfr
Environment=PYTHONPATH=/opt/nfr
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
